# supervisord/test.ini
[program:test]
command:php /Users/andrefigueira/Code/supervisord-manager/test/test.php
numprocs=1
autostart=true